# searxng-mcp.service
[Unit]
Description=SearxNG MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/workspace/mcp-searxng
Environment=SEARXNG_URL=http://172.17.0.1:2013
ExecStart=/workspace/mcp-searxng/.venv/bin/python mcp-searxng.py -b 0.0.0.0 -p 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
